# config/etc/supervisor.d/owncloud.ini
[group:owncloud]
programs=php-uwsgi,nginx

[program:php-uwsgi]
command = /usr/sbin/uwsgi --ini /etc/uwsgi/owncloud.ini
user = http
autostart = true

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true
